# src/deploy/dev_server_supervisor.ini

[program:redissvr]
directory = /data/
command=/data/redissvr/redissvr /data/redissvr/redissvr.toml
autostart=false
autorestart=true

[program:test]
directory = /data/
command=/data/test/test /data/test/test.toml
autostart=false
autorestart=true