# src/deploy/dev_server_supervisor.ini

[program:redissvr]
directory = /data/redissvr
command=/data/redissvr/redissvr /data/redissvr/redissvr.toml
autostart=false
autorestart=true

[program:discovery]
directory = /data/discovery
command=/data/discovery/discovery /data/discovery/discovery.toml
autostart=false
autorestart=true

[program:test]
directory = /data/test
command=/data/test/test /data/test/test.toml
autostart=false
autorestart=true